# cron/com.notiky.ernest.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Ernest morning brief via launchd. Runs the deterministic engine (no model, no
  connectors, no auth required), so it works even with no VPS and no open
  Claude session. Replace $HOME with your absolute home path before loading:
    launchctl load ~/Library/LaunchAgents/com.notiky.ernest.brief.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.notiky.ernest.brief</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>"$HOME/.ernest-cc/bin/ernest" brief</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>$HOME/.ernest-cc/logs/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>$HOME/.ernest-cc/logs/launchd.err</string>
</dict>
</plist>
